# crontab
10 2 * * * pipenv run python -u /home/jobrunner/all-tracks.py
